# skills/capability-evolver/com.evolver.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.evolver.weekly</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>echo "deprecated: do not load this plist; use OpenClaw cron job 903b61c4-51c7-4e09-ac1b-999ae1c2322c instead"</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/fx/.openclaw/workspace/skills/capability-evolver</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>0</integer>
        <key>Weekday</key>
        <integer>1</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/fx/.openclaw/workspace/logs/evolver-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/fx/.openclaw/workspace/logs/evolver-stderr.log</string>
</dict>
</plist>
