# scripts/com.openclaw.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.backup</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/fx/.openclaw/workspace/scripts/backup.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/fx/.openclaw/workspace</string>
    
    <key>StandardOutPath</key>
    <string>/Users/fx/.openclaw/workspace/logs/backup-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/fx/.openclaw/workspace/logs/backup-stderr.log</string>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
