# launchd/com.openclaw.cleanup-bailian-profile-cache.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.cleanup-browser-profile-cache</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/fx/.openclaw/workspace/scripts/cleanup-bailian-profile-cache.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/fx/.openclaw/logs/cleanup-browser-profile-cache.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/fx/.openclaw/logs/cleanup-browser-profile-cache.stderr.log</string>
</dict>
</plist>
